# files/etc/hotplug.d/usb/99-opentether
#!/bin/sh

PIDFILE=/tmp/opentether-hotplug.pid
LOCKDIR=/tmp/opentether.lock
STATEFILE=/tmp/opentether-device
PORT="$(uci -q get opentether.socks5.port || echo 1088)"

log() { logger -t opentether "$*"; }

forward_active() {
    local dev="$1" port="$2"
    adb -s "$dev" forward --list 2>/dev/null | awk -v p="tcp:${port}" '
        $2 == p && $3 == p { found=1; exit }
        END { exit(found ? 0 : 1) }'
}

stop_tunnel() {
    /etc/init.d/opentether stop >/dev/null 2>&1 || true
}

cleanup() {
    rm -f "$PIDFILE" "$STATEFILE"
    rm -rf "$LOCKDIR"
}

if [ "$ACTION" = "remove" ]; then
    stop_tunnel
    cleanup
    log "USB removed — tunnel stopped"
    exit 0
fi

[ "$ACTION" = "add" ] || exit 0
mkdir "$LOCKDIR" 2>/dev/null || exit 0

(
echo $$ > "$PIDFILE"
trap 'cleanup; exit 0' INT TERM EXIT

log "USB connected — waiting for ADB authorization"
adb devices >/dev/null 2>&1 || true

retries=90
while [ "$retries" -gt 0 ]; do
    DEVICE="$(adb devices 2>/dev/null | awk '/\tdevice$/{print $1; exit}')"
    if [ -n "$DEVICE" ]; then
        echo "$DEVICE" > "$STATEFILE"
        log "ADB authorized: $DEVICE"
        sleep 2

        adb -s "$DEVICE" forward tcp:${PORT} tcp:${PORT} >/dev/null 2>&1 || true

        if ! forward_active "$DEVICE" "$PORT"; then
            log "Forward failed"
            stop_tunnel
            exit 1
        fi

        if ! pgrep -f hev-socks5-tunnel >/dev/null; then
            log "Starting tunnel"
            /etc/init.d/opentether start
        fi

        log "Tunnel active"
        rm -rf "$LOCKDIR"

        while true; do
            sleep 5
            STATE="$(adb -s "$DEVICE" get-state 2>/dev/null || true)"

            if [ "$STATE" != "device" ]; then
                log "Device disconnected"
                stop_tunnel
                exit 0
            fi

            if ! forward_active "$DEVICE" "$PORT"; then
                log "Forward lost — restoring"
                adb -s "$DEVICE" forward tcp:${PORT} tcp:${PORT}
            fi

            if ! pgrep -f hev-socks5-tunnel >/dev/null; then
                log "Tunnel stopped — restarting"
                /etc/init.d/opentether start
            fi
        done
    fi

    sleep 1
    retries=$((retries - 1))
done

log "ADB authorization timeout"
stop_tunnel
exit 1
) &
